# deploy/ai-scanner-api.service
[Unit]
Description=AI Search Scanner — Next.js App
After=network.target

[Service]
Type=simple
User=jens
WorkingDirectory=/mnt/storage/aiscanner
Environment=PATH=/home/jens/.nvm/versions/node/v22.22.0/bin:/usr/bin:/usr/local/bin
EnvironmentFile=/mnt/storage/aiscanner/.env
Environment=PORT=8010
Environment=HOSTNAME=127.0.0.1
ExecStart=/home/jens/.nvm/versions/node/v22.22.0/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
